# supervisord.conf
[supervisord]
nodaemon=true

[program:nodejs]
command=node server.js
autorestart=true

[program:brackets]
command=brackets -p 8090 -j /workspace
autorestart=true

